# package/debian/archipel-core/archipel.slice
[Unit]
Description=Archipel delay tolerant network

[Slice]

[Install]
WantedBy=multi-user.target